# linux/rc.local
#!/bin/bash

if sha256sum /sys/devices/virtual/dmi/id/product_serial | cut -d " " -f 1 | grep b032ba4cc3803cb83af3dafb904d74c5e39aba4cc9696a2c89a0b95797902c59 ; then
  # not needed to most person.
  echo "0.0.0.0 news.163.com" >> /etc/hosts
  echo "0000:03:00.2"> /sys/bus/pci/drivers/xhci_hcd/unbind
  echo "0000:03:00.4"> /sys/bus/pci/drivers/xhci_hcd/unbind
fi
